# voice-chat/systemd/kavi-whisper-server.service
[Unit]
Description=Kavi - whisper-server (persistent STT for voice assistant)
After=graphical-session.target

[Service]
Type=simple
EnvironmentFile=%h/.config/kavi/kavi.env
ExecStart=%h/learn/Code/voice-assistant/voice-chat/scripts/start-whisper-server.sh
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
